# data/igord-event-publisher.service
[Unit]
Description=Igor event publisher
Requires=igord.service

[Service]
ExecStart=/usr/bin/igord-event-publisher

[Install]
WantedBy=default.target
