# data/igord.service
[Unit]
Description=igor daemon for automated testing
Wants=redis.service

[Service]
ExecStart=/usr/bin/igord

[Install]
WantedBy=default.target
